# scripts/com.policydb.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.policydb.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/grantgreeson/Documents/policydb/scripts/policydb-backup.sh</string>
  </array>

  <!-- Run daily at 2:00 AM -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- If the machine was asleep at 2 AM, run as soon as it wakes -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/grantgreeson/.policydb/backups/backup.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/grantgreeson/.policydb/backups/backup.log</string>

</dict>
</plist>
